# app/crontab
0 * * * * sh /code/entrypoint.sh && /code/update_shared_data.sh